# deploy/supervisor/clawdi-cloud.conf
# Supervisor config for clawdi-cloud backend on redpill.
#
# Install / refresh:
#   sudo cp deploy/supervisor/clawdi-cloud.conf /etc/supervisor/conf.d/
#   sudo supervisorctl reread
#   sudo supervisorctl update
#
# Daily ops:
#   sudo supervisorctl status clawdi-cloud-backend:*
#   sudo supervisorctl restart clawdi-cloud-backend:*
#   sudo tail -f /var/log/supervisor/clawdi-cloud-backend.log
#
# Two workers on 8070/8071 — nginx (cloud-api.clawdi.ai) round-robins across
# them via upstream. No celery/scheduler needed — clawdi-cloud has no
# background jobs (see docs/architecture.md).

[program:clawdi-cloud-backend]
directory=/opt/clawdi-cloud/backend
environment=PATH="/home/phala/.nvm/versions/node/v20.19.5/bin:/home/phala/.bun/bin:%(ENV_PATH)s"
command=/home/phala/.local/bin/uv run uvicorn app.main:app --port %(process_num)s
process_name=%(program_name)s-%(process_num)02d
numprocs=2
numprocs_start=8070
autorestart=true
redirect_stderr=true
killasgroup=true
stopasgroup=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=0
logfile_maxbytes=0
user=phala
